# deploy/slik-mobile-healthcheck.service
[Unit]
Description=SLIK Mobile Telegram Bot healthcheck watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/SLIK-Mobile
EnvironmentFile=/opt/SLIK-Mobile/.env
ExecStart=/opt/SLIK-Mobile/venv/bin/python /opt/SLIK-Mobile/bot/bot_healthcheck.py
